# cron.d/meister-db-backup
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""
0 3 * * * root /srv/meister/deploy/backup_mysql.sh >>/var/log/meister-backup.log 2>&1
